# .com.cowork.gitpush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Cowork git auto-commit LaunchAgent.
  Do NOT edit this file directly — run .cowork-install.sh which
  generates the real plist with hardcoded paths in ~/Library/LaunchAgents/.
  This copy is a template: __HOME__ gets replaced during install.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cowork.gitpush</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__HOME__/Projects/ClaudeCowork/.cowork-watcher.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__HOME__/Projects/ClaudeCowork</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/cowork-gitpush.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cowork-gitpush-error.log</string>
</dict>
</plist>
